# static/docker/laravel-horizon/supervisord.d/laravel-horizon.conf
[program:laravel-horizon]
user=horizon
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan horizon
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0